# assets/templates/etc/cron.d/whmcs
0 0 * * * root /usr/bin/php -q /var/www/public/crons/cron.php

